# runtime/deploy/systemd/alice-runtime.service
[Unit]
Description=Alice Runtime
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=alice-runtime
Group=alice-runtime
WorkingDirectory=/srv/alice-telegram-bot/runtime
Environment=NODE_ENV=production
Environment=ALICE_STATE_DIR=/var/lib/alice-runtime
Environment=ALICE_HOME=/var/lib/alice-runtime/home
Environment=ALICE_ENGINE_PORT=3380
Environment=ALICE_ERROR_LOG_PATH=/var/log/alice-runtime/alice-errors.log
Environment=ALICE_SANDBOX_RUNTIME=runsc
Environment=ALICE_CONTAINER_DEFAULT_PROFILE=sandboxed
EnvironmentFile=-/srv/alice-telegram-bot/runtime/.env
ExecStart=/usr/bin/node --import tsx --env-file=.env src/index.ts
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=mixed
UMask=0077
LimitNOFILE=65536

StateDirectory=alice-runtime
LogsDirectory=alice-runtime

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectClock=true
LockPersonality=true
RestrictNamespaces=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/srv/alice-telegram-bot/runtime/skills/store
ReadWritePaths=/srv/alice-telegram-bot/runtime/skills/system-bin
ReadWritePaths=/srv/alice-telegram-bot/runtime/skills/man
ReadWritePaths=/srv/alice-telegram-bot/runtime/skills/registry.json

[Install]
WantedBy=multi-user.target
